# supervisor/app.conf
[program:app]
directory=/home/ubuntu/apps/browsertap.com
command=/usr/local/bin/npm run app