# supervisor/website.conf
[program:website]
directory=/home/ubuntu/browsertap.com
command=/usr/local/bin/npm run website